# systemd/goal-agent.service
[Unit]
Description=Goal Agent API
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/goal-agent
EnvironmentFile=/home/pi/goal-agent/.env
ExecStart=/home/pi/goal-agent/.venv/bin/uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 1 \
    --loop uvloop \
    --access-log \
    --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goal-agent

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
